# etc/xdg/autostart/inpas-gui.desktop
[Desktop Entry]
Name=INPAS GUI
Exec=/opt/inpas-gui/notifier3.sh
Terminal=false
Type=Application
